# conf/etc/supervisor/conf.d/database.conf
[program:database-psql]
command=/opt/cloudcost/postgres.sh
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/postgres-stdout.log
stdout_logfile_maxbytes=10MB